# moonlight-voice/rootfs/etc/services.d/moonlight-voice/run
#!/usr/bin/with-contenv bash
set -euo pipefail

echo "[moonlight-voice] s6 service wrapper starting at $(date -Is)"
exec /opt/moonlight-voice/run.sh
